# cache_delete.conf
[program:cache_delete]
directory=/opt/cache_server
command=/opt/cache_server/venv/bin/python3 /opt/cache_server/cache_delete.py
autorestart=true
autostart=true
