# cdn_exporter.conf
[program:cdn_exporter]
directory=/opt/cache_server
command=/opt/cache_server/venv/bin/python3 /opt/cache_server/parse_log.py
autorestart=true
autostart=true
