# log_exporter.conf
[program:cdn_exporter]
directory=/home/nguyenquangminh/cache_server
command=/home/nguyenquangminh/venv/bin/python3 /home/nguyenquangminh/cache_server/parse_log.py
autorestart=true
autostart=true
